# backend/deploy/bloodbank_amazon_linux.service
[Unit]
Description=BloodBridge Flask Application
After=network.target

[Service]
Type=notify
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/bloodbank/backend
Environment="PATH=/home/ec2-user/bloodbank/backend/venv/bin"
Environment="FLASK_ENV=production"
EnvironmentFile=/home/ec2-user/bloodbank/backend/.env
ExecStart=/home/ec2-user/bloodbank/backend/venv/bin/gunicorn \
    --config /home/ec2-user/bloodbank/backend/deploy/gunicorn_config.py \
    app:app

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bloodbank

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
